# system/backends/void/runit/sol-zram/run
#!/bin/sh
set -eu
exec /usr/local/bin/apply-zram-policy.sh
