# system/backends/void/runit/sol-netd/run
#!/bin/sh
set -eu
exec /usr/local/bin/sol-netd
