# system/alpine/openrc/seatd
#!/sbin/openrc-run

description="Seat management daemon"

command="/usr/bin/seatd"
command_args="-u root"
command_user="root"

depend() {
    need localmount
}