# system/alpine/openrc/sol-zram
#!/sbin/openrc-run

description="Soliloquy zram swap"

command="/usr/local/bin/apply-zram-policy.sh"
command_background="no"

depend() {
    need localmount sol-kernel-policy
    before sold sol-session
}
