# system/alpine/openrc/sol-netd
#!/sbin/openrc-run

description="Soliloquy network state daemon"

command="/usr/local/bin/sol-netd"
command_background="yes"
pidfile="/run/${RC_SVCNAME}.pid"
export SOLILOQUY_NETD_STATE_JSON="${SOLILOQUY_NETD_STATE_JSON:-/run/soliloquy/netd/interfaces.json}"
export SOLILOQUY_NETD_RUNTIME_ENV="${SOLILOQUY_NETD_RUNTIME_ENV:-/run/soliloquy/netd/runtime-state.env}"
export SOLILOQUY_NETD_INTERVAL_SECS="${SOLILOQUY_NETD_INTERVAL_SECS:-5}"
export SOLILOQUY_RUNTIME_EVENTS_FILE="${SOLILOQUY_RUNTIME_EVENTS_FILE:-/run/soliloquy/runtime-events.log}"

now_unix_ms() {
    printf '%s000\n' "$(date +%s)"
}

record_runtime_event() {
    kind="$1"
    source="$2"
    detail="${3:-}"
    dir="$(dirname "${SOLILOQUY_RUNTIME_EVENTS_FILE}")"
    mkdir -p "${dir}"
    printf '%s\t%s\t%s\t%s\n' "$(now_unix_ms)" "${kind}" "${source}" "${detail}" >>"${SOLILOQUY_RUNTIME_EVENTS_FILE}" 2>/dev/null || true
}

start_pre() {
    record_runtime_event network_daemon_start openrc sol-netd
}

start_post() {
    record_runtime_event network_daemon_started openrc sol-netd
}

depend() {
    need localmount networking sol-kernel-policy
    before sold sol-session
}
