# system/alpine/openrc/sol-session
#!/sbin/openrc-run

description="Soliloquy session launcher"

command="/usr/local/bin/sol-session"
command_user="root"

depend() {
    need sold
}

respawn=YES
respawn_delay=5
