# system/alpine/openrc/sol-session
#!/sbin/openrc-run

description="Soliloquy session launcher"

command="/usr/local/bin/sol-session-start"

depend() {
    need sold seatd sol-kernel-policy sol-zram
}

respawn=YES
respawn_delay=5
