# system/alpine/openrc/sold
#!/sbin/openrc-run

description="Soliloquy local server"

command="/usr/local/bin/sold"
command_user="sold"
command_background="yes"
pidfile="/run/${RC_SVCNAME}.pid"
export SOLD_BIND="127.0.0.1:8080"
export SOLILOQUY_SYSTEM_CONFIG="/etc/soliloquy/system.json"
export SOLILOQUY_PLUGIN_STATE_FILE="/var/lib/soliloquy/system/plugin-state.json"

depend() {
    need localmount
    after seatd
}
