# system/alpine/rootfs-overlay/etc/profile.d/soliloquy.sh
export XDG_RUNTIME_DIR="/run/user/$(id -u)"
export SOL_TOKEN="${SOL_TOKEN:-dev-token-change-me}"
export SOL_UI_URL="${SOL_UI_URL:-file:///opt/sol/ui/index.html}"
export SOLILOQUY_SYSTEM_CONFIG="${SOLILOQUY_SYSTEM_CONFIG:-/etc/soliloquy/system.json}"
